# deploy/repager-egress-firewall.service
[Unit]
Description=RePager Docker egress protection
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/repager-egress-firewall
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
